# deploy/systemd/whaviso-deploy-host.timer
# whaviso: roda o atualizador do DEPLOY_HOST no boot e a cada 1 min. Instalar em
# /etc/systemd/system/whaviso-deploy-host.timer e habilitar:
#   systemctl enable --now whaviso-deploy-host.timer
[Unit]
Description=whaviso: checa o IP do VPS e atualiza o DEPLOY_HOST periodicamente

[Timer]
# Logo apos o boot (IP costuma mudar em reboot/remanejo) e depois a cada 1 min: encurta a
# janela de corrida em que um deploy pode disparar entre a troca de IP e o proximo tick
# (o secret ficaria velho e o SSH do CI daria timeout). O script e barato e idempotente.
OnBootSec=1min
OnUnitActiveSec=1min
# Se o VPS estava desligado na hora agendada, roda assim que voltar.
Persistent=true

[Install]
WantedBy=timers.target
