# deploy/systemd/whaviso-zap.service
# whaviso zap (scheduler + WhatsApp via Meta Cloud API, :3002).
# Instalar em /etc/systemd/system/whaviso-zap.service
#
# IMPORTANTE: instancia UNICA. O zap fala com a Meta Cloud API por HTTP (sem socket,
# sem QR, sem lock de sessao) e recebe o inbound por webhook. Ainda assim rode UMA
# so copia: o scheduler drena a outbox e nao deve competir consigo mesmo.
[Unit]
Description=whaviso zap (scheduler + WhatsApp via Meta Cloud API :3002)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=whaviso
Group=whaviso
WorkingDirectory=/opt/whaviso/app/backend/apps/zap
EnvironmentFile=/etc/whaviso/whaviso.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/env npm run start
Restart=on-failure
RestartSec=10

# Area de escrita de runtime do zap, fora do checkout (que fica imutavel). Persiste
# entre deploys e reboots. Nao apague este diretorio.
ReadWritePaths=/var/lib/whaviso

# ── Hardening (systemd sandbox) ──
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryHigh=512M
MemoryMax=768M

[Install]
WantedBy=multi-user.target
