# deploy/systemd/whaviso-zap.service
# whaviso zap (scheduler + WhatsApp/Baileys, :3002).
# Instalar em /etc/systemd/system/whaviso-zap.service
#
# IMPORTANTE: instancia UNICA. O zap segura o socket do WhatsApp e um lock de
# arquivo (.baileys.lock). Nunca rode duas copias: dois processos na mesma sessao
# causam o erro 440. Por isso nada de escala horizontal aqui.
[Unit]
Description=whaviso zap (scheduler + WhatsApp :3002)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=whaviso
Group=whaviso
WorkingDirectory=/opt/whaviso/app/backend/apps/zap
EnvironmentFile=/etc/whaviso/whaviso.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/env npm run start
Restart=on-failure
RestartSec=10

# A sessao do Baileys (e o lock e o qr.png) vivem aqui e PRECISAM persistir
# entre deploys e reboots. Nao apague este diretorio.
ReadWritePaths=/var/lib/whaviso

# ── Hardening (systemd sandbox) ──
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryHigh=512M
MemoryMax=768M

[Install]
WantedBy=multi-user.target
